# Wolf_elementary OS_0.3/etc/skel/.config/autostart/custom-command0.desktop
[Desktop Entry]
Name[en]=Custom Command
Comment[en]=xbindkeys
Exec=xbindkeys
Icon=application-default-icon
X-GNOME-Autostart-enabled=true
Type=Application
