# files/etc/cron.daily/distroshare-updater
#!/bin/sh

/usr/sbin/distroshare-updater
