# Ubuntu_15.04/etc/udev/rules.d/99-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{HDMI}=="1", RUN+="/usr/local/bin/hdmi-change.sh"
